# org.tizen.quickpanel.manifest
<manifest>
	<define>
	    <domain name="org.tizen.quickpanel"/>
	</define>
	<assign>
		<filesystem path="/etc/init.d/quickpanel" label="_" exec_label="none"/>
	</assign>
	<request>
	    <domain name="org.tizen.quickpanel"/>
	</request>
</manifest>
